# dyn-update.service.in
# Copyright (c) 2021 Garry T. Williams

[Unit]
Description=Keep Dyn DNS Updated
Documentation=man:dyn-update(1)

[Service]
Type=exec
Restart=always
RestartSec=10
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=-@HOME@/.cache/dyn
Nice=19
ExecStart=@HOME@/.local/bin/dyn-update

[Install]
WantedBy=default.target

# vim: set sw=4 sts=4 ts=8 et syn=systemd:
