# dyn-update.service.in
# Copyright (c) 2021 Garry T. Williams
#
# This is free software.  You can redistribute it and/or modify it
# under the terms of the GNU General Public License version 3 or any
# later version.  You should have a copy of the GNU General Public
# License supplied with this program.

[Unit]
Description=Keep Dyn DNS Updated
Documentation=man:dyn-update(1)

[Service]
Type=exec
Restart=always
RestartSec=10
NoNewPrivileges=true
Nice=19
ExecStart=@HOME@/.local/bin/dyn-update

[Install]
WantedBy=default.target

# vim: set sw=4 sts=4 ts=8 et syn=systemd:
